# .yuri_misc_sys_files/updates-available.service
[Unit]
Description=Check for Arch updates

[Service]
Type=simple
ExecStart=/usr/local/bin/updates-available.sh
